# pkg/baraud-base/baraud-base.install
post_install(){
	post_upgrade
}

post_upgrade(){
	systemctl enable dhcpcd
	sed -i '1 i\Include $XDG_CONFIG_HOME/ssh/config' /etc/ssh/ssh_config
}
