# pkg/baraud-desktop/baraud-desktop.install
post_install(){
	post_upgrade
}

post_upgrade(){
	cp /usr/share/baraud/lightdm/* /etc/lightdm/
	systemctl enable lightdm
}
